# services/aria2.service
[Unit]
Description=Aria2c download tool

[Service]
ExecStart=aria2c

[Install]
WantedBy=default.target
